# scripts/com.ai-agent-history-rag.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ai-agent-history-rag.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USERNAME/.local/bin/uv</string>
        <string>--directory</string>
        <string>/path/to/ai-agent-history-rag-mcp</string>
        <string>run</string>
        <string>ai-agent-history-rag-daemon</string>
        <string>start</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAUDE_HISTORY_RAG_OLLAMA_HOST</key>
        <string>http://YOUR_OLLAMA_HOST:11434</string>
        <key>CLAUDE_HISTORY_RAG_OLLAMA_MODEL</key>
        <string>bge-m3</string>
        <key>CLAUDE_HISTORY_RAG_BATCH_SIZE</key>
        <string>64</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.claude-history-rag/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.claude-history-rag/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/path/to/ai-agent-history-rag-mcp</string>
</dict>
</plist>
